# deploy/systemd/nexaplay-deploy.service
[Unit]
Description=NexaPlay automatic deployment
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/nexaplay
ExecStart=/opt/nexaplay/deploy/auto-deploy.sh
TimeoutStartSec=30min
Nice=10

