# deploy/systemd/nexaplay-deploy.timer
[Unit]
Description=Check NexaPlay GitHub updates every two minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
